# debian/ntpdate.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]; then
	if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1:4.2.2+dfsg.2-3; then
		rm -f /etc/init.d/ntpdate
		update-rc.d ntpdate remove >/dev/null
	fi
fi

#DEBHELPER#
